# .docker/config/crontab/default
* * * * * cd /var/www && /usr/local/bin/php index.php >> /var/log/crontab.log 2>&1
